# debian/postinst
#!/bin/sh
# postinst script for zswapd.
#
# See: dh_installdeb(1).

set -e

# Summary of how this script can be called:
#        * <postinst> 'configure' <most-recently-configured-version>
#        * <old-postinst> 'abort-upgrade' <new version>
#        * <conflictor's-postinst> 'abort-remove' 'in-favour' <package>
#          <new-version>
#        * <postinst> 'abort-remove'
#        * <deconfigured's-postinst> 'abort-deconfigure' 'in-favour'
#          <failed-install-package> <version> 'removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package.

case "${1}" in
    configure|upgrade|reinstall|install             )
        # Restart the systemd-modules-load service to ensure the zram module is loaded
        systemctl restart systemd-modules-load.service

        # Check if zswapd is already enabled if so reload the daemon
        case "$(systemctl is-enabled zswapd.service)" in
            enabled             ) systemctl daemon-reload ; systemctl restart zswapd.service ;;
            disabled|inactive   ) systemctl enable --now zswapd.service                      ;;
            *                   )
                echo "zswapd service is in an unknown state, please check manually." >&2
            ;;
        esac

    ;;
    *                                               )
        echo "postinst called with unknown argument: ${1}" >&2
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
